# template/etc/supervisord.conf.in
[unix_http_server]
file=/tmp/supervisor.sock
 
[supervisord]
logfile=/Users/seang/code/r-tree-server/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=warn
pidfile=/Users/seang/code/r-tree-server/var/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=nobody
childlogdir=/Users/seang/code/r-tree-server/log/
 
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
 
[group:myapp]
programs=vt-serve
 
[program:vt-serve]
command=/Users/seang/code/r-tree-server/scripts/vt-serve -c /Users/seang/code/r-tree-server/etc/vaytrou.conf
directory=/Users/seang/code/r-tree-server/
autorestart=true
redirect_stderr=true
stdout_logfile=/Users/seang/code/r-tree-server/log/vt-serve.log
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=warn
stopsignal=INT 
