# template/etc/supervisord.conf.in
[unix_http_server]
file=/tmp/supervisor.sock
 
[supervisord]
logfile=${home}/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=warn
pidfile=i${home}/var/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=nobody
childlogdir=${home}/log/supervisord/
 
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
 
[group:myapp]
programs=server1
 
[program:server1]
command=${home}/scripts/vt-serve -c ${home}/etc/vaytrou.conf
directory=${home}
autorestart=true
redirect_stderr=true
stdout_logfile=${home}/log/server1.log
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=warn
stopsignal=INT 
